# launchd/com.roncompton.avcontrol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.roncompton.avcontrol</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/roncompton/AVCONTROL/av_proxy.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/roncompton/AVCONTROL</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/roncompton/AVCONTROL/logs/avcontrol.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/roncompton/AVCONTROL/logs/avcontrol.err.log</string>
</dict>
</plist>
